# dockerfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:main]
command=node /app/core/server.js --redis

[program:matdaemon]
command=redis-server
