# dockerfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:main]
directory=/app/image
command=node core/server.js --redis

[program:redis]
command=redis-server

[program:nginx]
command=nginx